# profile.d/tty.sh
#!/bin/sh
[ "$(tty)" = "/dev/tty1" ] && startx &>/dev/null
